# hellosunshine.conf
description "Hello world Deamon"

# Start when the system is ready to do networking.
start on started elastic-network-interfaces

# Stop when the system is on its way down.
stop on shutdown

respawn
script
    exec su --session-command="https://github.com/Zand76/SWDV660/blob/main/hellosunshine.js" ec2-user
end script